# n8n-auto-heal.service
[Unit]
Description=n8n Auto-Heal HTTP Listener
Documentation=https://github.com/valtterimelkko/n8n-http-listener
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/n8n-http-listener

# Environment
Environment="PATH=/root/.local/bin:/usr/bin:/bin"
Environment="GMAIL_EMAIL=valtteri.melkko@gmail.com"
Environment="GMAIL_APP_PASSWORD=xxorioyqoqlcmgaj"
Environment="HOME=/root"

# Main executable
ExecStart=/usr/bin/python3 /root/n8n-http-listener/listener.py

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=n8n-auto-heal

# Security (run as root is required for kimi CLI access)
# If you want to run as non-root, adjust permissions accordingly

[Install]
WantedBy=multi-user.target
